# images/mkosi.extra/usr/lib/systemd/system/getty@tty1.service.d/90-rdi-installer.conf
[Service]
ExecStart=
ExecStart=-agetty --chdir /root -l /usr/bin/rdi-installer --autologin root --noclear - $TERM
StandardInput=tty
StandardOutput=tty
